# src/node-nginx/24/container/templates/supervisor/conf.d/node.service.mode-web.conf
[program:node-server]
command=gosu www-data ${NODE_WEB_COMMAND}
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
startsecs=5
startretries=3
