# src/python-nginx/3.13/container/templates/supervisor/conf.d/gunicorn.service.mode-web.conf
[program:gunicorn]
command=/opt/venv/bin/gunicorn -c /etc/gunicorn.conf.py ${PYTHON_APP_MODULE}
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
